# etc/systemd/system/docker.service.d/90-d26a-network.conf
[Service]
Environment="DOCKER_OPTS=--bridge=d26a --ip-masq=false"
